# imap-notify.service
[Unit]
Description=IMAP notifier

[Service]
ExecStart=/usr/bin/imap-notify
Restart=always

[Install]
WantedBy=default.target
